# deploy/launchd/com.example.ngrokd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.example.ngrokd</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/ngrokd</string>
    <string>-domain=example.com</string>
    <string>-tlsCrt=/usr/local/etc/ngrok/tls/example.com.crt</string>
    <string>-tlsKey=/usr/local/etc/ngrok/tls/example.com.key</string>
    <string>-authToken=change-me-long-random-token</string>
    <string>-httpAddr=127.0.0.1:8080</string>
    <string>-httpsAddr=127.0.0.1:8443</string>
    <string>-tunnelAddr=0.0.0.0:4443</string>
    <string>-log=stdout</string>
    <string>-log-level=INFO</string>
    <string>-maxConnections=1024</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/usr/local/var/log/ngrokd.log</string>
  <key>StandardErrorPath</key>
  <string>/usr/local/var/log/ngrokd.err.log</string>
  <key>WorkingDirectory</key>
  <string>/usr/local/var/lib/ngrok</string>
</dict>
</plist>
