# deploy/launchd/com.example.ngrok-client.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.example.ngrok-client</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/ngrok</string>
    <string>-config=/usr/local/etc/ngrok/client.yml</string>
    <string>-log=stdout</string>
    <string>start-all</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/usr/local/var/log/ngrok-client.log</string>
  <key>StandardErrorPath</key>
  <string>/usr/local/var/log/ngrok-client.err.log</string>
  <key>WorkingDirectory</key>
  <string>/usr/local/var/lib/ngrok</string>
</dict>
</plist>
